# tools/ansible/include/systemd_cron/templates/job.service
[Unit]
Description={{ name }}

[Service]
Type=oneshot
ExecStart=/etc/systemd/cron-scripts/{{ name }}.sh
