# tools/ansible/include/systemd_cron/templates/job.timer
[Unit]
Description=Runs {{ name }}.service

[Timer]
OnCalendar={{ time }}
Unit={{ name }}.service

[Install]
WantedBy=multi-user.target
